# etc/NetworkManager/dispatcher.d/update-i3blocks.sh
#!/bin/sh

/home/saumon/scripts/send-i3blocks-signal.sh 14
sleep 1
/home/saumon/scripts/send-i3blocks-signal.sh 14
